# deploy/systemd/promptdeck-api.service
[Unit]
Description=promptDeck FastAPI (uvicorn)
After=network.target postgresql.service

[Service]
Type=simple
User=promptdeck
Group=promptdeck
WorkingDirectory=/opt/promptDeck/backend
EnvironmentFile=/etc/promptdeck/api.env
ExecStart=/opt/promptDeck/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8005 --proxy-headers
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
